# systemd/rv-media-player.service
[Unit]
Description=RV Media Player - Jellyfin Offline Media Server
Documentation=https://github.com/your-repo/rv-media-player
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# The service runs as the desktop user so VLC can open windows on their
# Wayland/X11 session.  install.sh sets this automatically from $SUDO_USER.
# If editing by hand, set User/Group to your login username.
User=riley
Group=riley
SupplementaryGroups=audio video dialout plugdev
WorkingDirectory=/opt/rv-media-player
Environment=PYTHONPATH=/opt/rv-media-player
Environment=FLASK_APP=app/app.py
Environment=FLASK_ENV=production
# systemd strips display session variables — inject them so VLC can reach the
# compositor when launched from this service.
# DISPLAY=:0 covers X11 and XWayland (Wayland compositors expose XWayland here).
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/riley/.Xauthority
# Wayland: Pi OS Bookworm (Wayfire) uses wayland-0.
Environment=WAYLAND_DISPLAY=wayland-0
# Must match the UID of User= above.  Run 'id -u riley' to confirm.
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStartPre=/bin/mkdir -p /opt/rv-media-player/logs /opt/rv-media-player/data
ExecStart=/opt/rv-media-player/venv/bin/python -m app.app
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rv-media-player

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
# false — service user needs ~/.Xauthority for X11/XWayland auth.
ProtectHome=false
ReadWritePaths=/opt/rv-media-player /media /tmp

LimitNOFILE=65536

# Process priority
Nice=-5
IOSchedulingClass=1
IOSchedulingPriority=4

[Install]
WantedBy=multi-user.target
